# contrib/miner-agent.slice
[Unit]
Description=Slice for miner-agent resource accounting

[Slice]
CPUAccounting=yes
MemoryAccounting=yes
TasksAccounting=yes
